# proj-broadcast.service
[Unit]
Description=Start Broadcast server
After=network.target

[Service]
ExecStart=/opt/utils/ProjectComm/broadcast.py

[Install]

